# scripts/monitoring/claude-cpu-monitor.timer
[Unit]
Description=Запускать CPU monitor каждую минуту

[Timer]
OnBootSec=60
OnUnitActiveSec=1min
AccuracySec=10

[Install]
WantedBy=timers.target
