# scripts/firewall/claude-egress-monitor.service
[Unit]
Description=Claude Bot Egress Monitor (per-container traffic accounting)
After=claude-firewall.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/claude-egress-monitor.sh
StandardOutput=journal
StandardError=journal
